# deploy/clipku-backup.service
[Unit]
Description=Clipku database and source backup

[Service]
Type=oneshot
User=ubuntu
ExecStart=/home/ubuntu/clipku-streaming/scripts/backup.sh
