# deploy/clipku-monitor.service
[Unit]
Description=Clipku playback audit and trending refresh
After=network-online.target mysql.service

[Service]
Type=oneshot
User=ubuntu
WorkingDirectory=/home/ubuntu/clipku-streaming
EnvironmentFile=/home/ubuntu/clipku-streaming/.env
ExecStart=/usr/bin/npm run audit:playback
ExecStart=/usr/bin/npm run trending:update
